# armsr/armv8/N1/files/etc/uci-defaults/95-fix-tar
#!/bin/sh
ln -sf /usr/bin/bsdtar /bin/tar
exit 0