# armsr/armv8/N1/files/etc/uci-defaults/99-fix-tar
#!/bin/sh
ln -sf /usr/bin/bsdtar /usr/bin/tar
exit 0